# supervisord.conf
[supervisord]
nodaemon=true

[program:mongo]
command=/usr/bin/mongod
priority=1

[program:node]
command=/usr/local/bin/node server.js
directory=/opt/loowid
priority=2
